# debian/nrpe-ng.service
[Unit]
Description=The next generation Nagios Remote Plugin Executor

[Service]
Type=forking
PIDFile=/run/nagios/nrpe-ng.pid
ExecStart=/usr/sbin/nrpe-ng --config /etc/nagios/nrpe-ng.cfg --daemon

[Install]
WantedBy=multi-user.target
